# services/disabled/geki.service
[Unit]
Description=Geki
After=multi-user.target
 
[Service]
Type=idle
ExecStart=/usr/bin/python3 /var/srvcs/services/scripts/game.py
Restart=Restart=on-failure
 
[Install]
WantedBy=multi-user.target